# webapp.path
[Unit]
Description= waiting for startup Script

[Path] 
PathExists = /opt/startupdb.sh

[Install] 
WantedBy = multi-user.target


